# systemd/anima-broker.service
[Unit]
Description=Anima Hardware Broker - Lumen's Body (Sensors & Display)
After=network.target redis.service anima-restore.service
Wants=network-online.target
# Broker must start before MCP server (mind)
Before=anima.service
# PartOf: when anima is restarted (e.g. deploy), broker also restarts
PartOf=anima.service
# Prevent rapid restarts - wait for hardware to stabilize
OnFailure=anima-broker-failed.service

[Service]
Type=simple
User=unitares-anima
Group=unitares-anima
WorkingDirectory=/home/unitares-anima/anima-mcp

# Secrets from ~/.anima/anima.env (see config/anima.env.example)
EnvironmentFile=-/home/unitares-anima/.anima/anima.env

# Environment variables
# ANIMA_ID is read from ~/anima.db; do not bake an operator's UUID into a
# public template. Set ANIMA_ID in ~/.anima/anima.env only when forcing a
# specific creature_id (e.g. restoring an identity).
# UNITARES_URL is operator-specific (Tailscale-resolved). Set it in
# ~/.anima/anima.env using the hostname form so it survives Mac IP changes.
Environment="ANIMA_DB=/home/unitares-anima/.anima/anima.db"
Environment="ANIMA_CONFIG=/home/unitares-anima/anima-mcp/anima_config.yaml"
Environment="ANIMA_GOVERNANCE_INTERVAL_SECONDS=180"
Environment="PYTHONUNBUFFERED=1"
Environment="PYTHONIOENCODING=utf-8"
Environment="LC_ALL=C.UTF-8"

# Quiet period after stop: let connections drain, give WiFi chip time to settle.
ExecStopPost=/bin/sleep 3

# Executable - runs anima-creature entry point (pyproject.toml → stable_creature:main)
# Add timeout to prevent hanging on hardware initialization
ExecStart=/home/unitares-anima/anima-mcp/.venv/bin/anima-creature
TimeoutStartSec=30
TimeoutStopSec=10

# Restart policy - broker should be very stable
# Use exponential backoff to prevent reboot loops
Restart=on-failure
RestartSec=10
# Increase delay between restarts to prevent loops
# After 3 failures in 5 minutes, wait 5 minutes before allowing restart
StartLimitInterval=600
StartLimitBurst=5
StartLimitAction=none

# Resource limits
LimitNOFILE=65536
MemoryMax=256M
CPUQuota=30%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=anima-broker

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=no
ReadWritePaths=/home/unitares-anima/anima-mcp /home/unitares-anima/.anima /dev/shm

[Install]
WantedBy=multi-user.target
