# systemd/anima.service
[Unit]
Description=Anima MCP Server - Lumen's Mind (MCP Interface)
After=network.target anima-broker.service anima-restore.service
Requires=anima-broker.service
Wants=network-online.target
# Broker owns sensors, writes to shared memory; server owns DB (no contention)

[Service]
Type=simple
User=unitares-anima
Group=unitares-anima
WorkingDirectory=/home/unitares-anima/anima-mcp

# Environment variables (secrets from ~/.anima/anima.env — see config/anima.env.example)
EnvironmentFile=-/home/unitares-anima/.anima/anima.env
# ANIMA_ID is read from ~/anima.db; do not bake an operator's UUID into a
# public template. Set ANIMA_ID in ~/.anima/anima.env only when forcing a
# specific creature_id (e.g. restoring an identity).
# UNITARES_URL is operator-specific (Tailscale-resolved). Set it in
# ~/.anima/anima.env using the hostname form so it survives Mac IP changes.
Environment="ANIMA_DB=/home/unitares-anima/.anima/anima.db"
Environment="ANIMA_CONFIG=/home/unitares-anima/anima-mcp/anima_config.yaml"
Environment="ANIMA_OAUTH_ISSUER_URL=https://lumen.cirwel.org"
Environment="ANIMA_OAUTH_AUTO_APPROVE=true"
Environment="PYTHONUNBUFFERED=1"
# Server reads broker-published sensor data from /dev/shm — never opens
# /dev/i2c-1. Pre-fix, lazy PiSensors() in the server held a duplicate
# I2C bus handle that contaminated BMP280 state (silent at 0x77 for ~42h
# on 2026-04-26 → 2026-04-28). Broker (anima-broker.service) is the sole
# I2C owner; it does NOT set this env var.
Environment="ANIMA_SENSORS_BACKEND=shm"

# Quiet period after stop: let existing connections drain before restart begins.
# Prevents the connection burst that crashes brcmfmac WiFi firmware on Pi.
ExecStopPost=/bin/sleep 3

# Kill any zombie process holding the port before starting
ExecStartPre=/bin/sh -c 'fuser -k 8766/tcp 2>/dev/null || true; sleep 2'

# Executable
ExecStart=/home/unitares-anima/anima-mcp/.venv/bin/anima --http --host 0.0.0.0 --port 8766

# Restart policy - allow more retries for undervoltage recovery
Restart=on-failure
RestartSec=20
StartLimitInterval=600
StartLimitBurst=5
# If 3 failures in 5 min, reboot Pi as last resort recovery
StartLimitAction=none

# Resource limits
LimitNOFILE=65536
MemoryMax=512M
CPUQuota=50%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=anima

# Security
# NoNewPrivileges omitted: server needs sudo for operational tasks
# (service sync, systemctl restart, system_service handler)
PrivateTmp=true
ProtectSystem=strict
ProtectHome=no
ReadWritePaths=/home/unitares-anima/anima-mcp /home/unitares-anima/.anima /dev/shm

[Install]
WantedBy=multi-user.target
