# anima_broker/systemd/anima-broker-ex.service
[Unit]
# Phase 0 Elixir broker. Runs ALONGSIDE the Python anima-broker.service and
# writes to a SHADOW SHM path (ANIMA_SHM_PATH) so it cannot affect the live
# creature. Do not point ANIMA_SHM_PATH at /dev/shm/anima_state.json until the
# Phase 1 sensor cutover.
Description=Anima Elixir Broker (Phase 0, shadow SHM)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=unitares-anima
WorkingDirectory=/home/unitares-anima/anima-mcp/anima_broker
Environment=ANIMA_SHM_PATH=/dev/shm/anima_state.shadow.json
Environment=ANIMA_TICK_MS=2000
# Build first with: MIX_ENV=prod mix release
ExecStart=/home/unitares-anima/anima-mcp/anima_broker/_build/prod/rel/anima_broker/bin/anima_broker start
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
