# config/anima-creature.service
[Unit]
Description=Anima Creature - Lumen's Hardware Broker
After=network.target
Before=anima.service

[Service]
Type=simple
User=unitares-anima
Group=unitares-anima
WorkingDirectory=/home/unitares-anima/anima-mcp

# Secrets from ~/.anima/anima.env (see config/anima.env.example)
EnvironmentFile=-/home/unitares-anima/.anima/anima.env

# ANIMA_ID is normally read from ~/anima.db; only set here as an override.
# See anima.env.example for the env-var form. A fresh install (no DB) mints
# a new creature_id automatically — do NOT bake the original operator's UUID
# into a public template.
Environment="ANIMA_DB=/home/unitares-anima/anima.db"
Environment="ANIMA_CONFIG=/home/unitares-anima/anima-mcp/anima_config.yaml"
# Set UNITARES_URL in ~/.anima/anima.env (operator-specific tailnet; see anima.env.example)
Environment="PYTHONUNBUFFERED=1"

ExecStart=/home/unitares-anima/anima-mcp/.venv/bin/anima-creature

Restart=on-failure
RestartSec=10
StartLimitInterval=300
StartLimitBurst=5

LimitNOFILE=65536
MemoryMax=256M

StandardOutput=journal
StandardError=journal
SyslogIdentifier=anima-creature

[Install]
WantedBy=multi-user.target
